# ubuntu/config/templates/dropbox.service.j2
[Unit]
Description=dropbox


[Service]
Type=simple
ExecStart=/home/{{ user }}/.dropbox-dist/dropboxd
ExecStop=


[Install]
WantedBy=multi-user.target
